# lib/kdeploy/templates/app.service.erb
[Unit]
Description=<%= application %> application server
After=network.target
Requires=<%= application %>-worker.service

[Service]
Type=simple
User=<%= user %>
Group=<%= user %>
WorkingDirectory=<%= deploy_to %>
Environment=RAILS_ENV=<%= environment %>
Environment=RACK_ENV=<%= environment %>
Environment=NODE_ENV=<%= environment %>
Environment=PORT=<%= app_port %>
Environment=DATABASE_URL=<%= database_url %>
Environment=REDIS_URL=<%= redis_url %>
Environment=MEMCACHED_URL=<%= memcached_url %>

# Resource limits
LimitNOFILE=65536
LimitNPROC=65536
MemoryHigh=<%= max_memory || '4G' %>
MemoryMax=<%= max_memory || '4G' %>
TasksMax=infinity

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=0

# Process management
ExecStart=/usr/local/bin/bundle exec puma -C config/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=30

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target