# examples/templates/app.service.erb
# Systemd service file for <%= application %>
# Generated by Kdeploy on <%= Time.now.strftime('%Y-%m-%d %H:%M:%S') %>

[Unit]
Description=<%= description || "#{application} web application" %>
After=network.target
Wants=network.target

[Service]
Type=simple
User=<%= user %>
Group=<%= user %>
WorkingDirectory=<%= deploy_to %>

# Command to start the application
ExecStart=/usr/bin/<%= exec_start || 'node server.js' %>
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5

# Restart policy
Restart=always
RestartSec=10

# Environment variables
Environment=NODE_ENV=<%= environment || 'production' %>
Environment=PORT=<%= app_port || 3000 %>
<% if postgres_port %>
Environment=DB_PORT=<%= postgres_port %>
<% end %>

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=<%= deploy_to %>/logs <%= deploy_to %>/tmp <%= deploy_to %>/uploads

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Standard I/O
StandardOutput=journal
StandardError=inherit
SyslogIdentifier=<%= application %>

[Install]
WantedBy=multi-user.target